# init-scripts/dofler-client.conf
description "Dofler Parsers"
author "steve@chigeek.com"

start on runlevel [234]
stop on runlevel [0156]

exec /usr/bin/dofler-client

respawn