# systemd.nsfw.service
[Unit]
Description=Dofler NSFW Service

[Service]
Restart=Always
ExecStart=/usr/bin/docker run --name nsfw_service -p 4000:8080 -v /opt/dofler/images:/images nsfw:monolith
ExecStop=/usr/bin/docker stop nsfw_service

[Install]
WantedBy=multi-user.target